# conf/brewlog.2019.upstart.conf
description "uWSGI service for brewlog"

start on runlevel [2345]
stop on runlevel [!2345]

setuid jarek
setgid www-data

env PATH=/home/jarek/brewlog/_venv/bin
env BREWLOG_CONFIG_LOCAL=/home/jarek/brewlog/config/config_local.py
env BREWLOG_CONFIG_SECRETS=/home/jarek/brewlog/config/secrets.py
chdir /home/jarek/brewlog
exec uwsgi --ini config/brewlog.ini
